# apps/broker/broker.service
# Message Broker Service
#
# Internal service for Claude streaming - runs alongside Next.js
# Owns: stream state machines, concurrency limits, graceful shutdown
#
# Install:
#   sudo cp apps/broker/broker.service /etc/systemd/system/claude-broker.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable claude-broker
#   sudo systemctl start claude-broker
#
# Logs:
#   journalctl -u claude-broker -f

[Unit]
Description=Claude Message Broker
After=network.target

[Service]
Type=simple
User=root
Group=root

# Working directory
WorkingDirectory=/root/webalive/alive

# Start command - uses Bun for performance
ExecStart=/root/.bun/bin/bun apps/broker/src/server.ts

# Environment
Environment=NODE_ENV=production
Environment=BROKER_PORT=3001
Environment=BROKER_HOST=127.0.0.1

# Load secrets from environment file
EnvironmentFile=/root/webalive/alive/.env.broker

# Restart policy
Restart=always
RestartSec=5

# Graceful shutdown (30 seconds to drain streams)
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Security hardening (internal service)
NoNewPrivileges=false
PrivateTmp=true

# Resource limits
MemoryMax=2G
TasksMax=1000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-broker

[Install]
WantedBy=multi-user.target
