# ops/systemd/alive-dev.service
[Unit]
Description=Alive Dev
After=network.target

[Service]
Type=simple
WorkingDirectory=__ALIVE_ROOT__/apps/web
ExecStart=/usr/local/bin/bun run dev

EnvironmentFile=-__ALIVE_ROOT__/apps/web/.env.local

Environment="HOME=/root"
Environment="NODE_ENV=development"
Environment="PORT=8997"
Environment="BRIDGE_ENV=local"

MemoryMax=3G
MemoryHigh=2G
CPUQuota=200%
TasksMax=512
LimitNOFILE=4096

Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=alive-dev

[Install]
WantedBy=multi-user.target
