# ops/systemd/claude-bridge-dev.service
[Unit]
Description=Claude Bridge Development Server
Documentation=https://github.com/AliveBridgeNL/claude-bridge
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=/root/webalive/claude-bridge/apps/web

# Development server with Turbopack
ExecStart=/usr/local/bin/bun run dev

# Environment files (loaded in order, later ones override)
EnvironmentFile=-/root/webalive/claude-bridge/apps/web/.env.development

# Additional environment
Environment="NODE_ENV=development"
Environment="PORT=8997"

# Resource limits (more generous for dev)
MemoryMax=4G
MemoryHigh=3G
CPUQuota=400%
TasksMax=1024
LimitNOFILE=8192

# Auto-restart on failure
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-bridge-dev

[Install]
WantedBy=multi-user.target
