# ops/systemd/claude-bridge-staging.service
[Unit]
Description=Claude Bridge Staging
Documentation=https://github.com/AliveBridgeNL/claude-bridge
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=/root/webalive/claude-bridge/.builds/staging/current/standalone

# Main executable
ExecStart=/usr/local/bin/bun /root/webalive/claude-bridge/.builds/staging/current/standalone/apps/web/server.js

# Environment
Environment="NODE_ENV=production"
Environment="PORT=8998"
Environment="BRIDGE_API_PORT=8998"
EnvironmentFile=-/root/webalive/claude-bridge/apps/web/.env.staging

# Resource limits
MemoryMax=2G
MemoryHigh=1.5G
CPUQuota=200%
TasksMax=512
LimitNOFILE=4096

# Quick graceful shutdown - Bun/Next.js doesn't need 90s
TimeoutStopSec=10s

# Restart configuration
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
