# ops/systemd/fix-permissions.service
[Unit]
Description=Fix site directory permissions
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/webalive/alive/scripts/fix-permissions.sh
WorkingDirectory=/root/webalive/alive
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
