# ops/systemd/fix-permissions.service
[Unit]
Description=Fix site directory permissions
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=__ALIVE_ROOT__/scripts/fix-permissions.sh
WorkingDirectory=__ALIVE_ROOT__
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
