# ops/systemd/preview-proxy.service
[Unit]
Description=Alive Preview Proxy (Go)
Documentation=https://github.com/eenlars/alive
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/webalive/alive/apps/preview-proxy

ExecStart=/root/webalive/alive/apps/preview-proxy/dist/preview-proxy

Environment="PREVIEW_BASE=alive.best"
Environment="JWT_SECRET=L+Y259oGhg4I9ik/H+Tby7ybum+t1fFHripWTXQKa2w="
Environment="PORT_MAP_PATH=/var/lib/alive/generated/port-map.json"
Environment="LISTEN_ADDR=:5055"
Environment="FRAME_ANCESTORS=https://app.alive.best,https://staging.alive.best,https://dev.alive.best"

MemoryMax=128M
MemoryHigh=96M
CPUQuota=50%
TasksMax=128
LimitNOFILE=4096

Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s
TimeoutStartSec=10s

[Install]
WantedBy=multi-user.target
