# ops/systemd/website-backup.service
[Unit]
Description=Backup websites to GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/webalive/alive/scripts/backup-websites.sh
WorkingDirectory=/root/webalive/alive
Environment=PATH=/root/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
