# apps/mcp/stealth-request/stealth-request.service
[Unit]
Description=Stealth Request Web Scraping Service
After=network.target
Documentation=file:///root/webalive/alive/apps/mcp/stealth-request/README.md

[Service]
Type=simple
User=root
WorkingDirectory=/root/webalive/alive
ExecStart=/root/.bun/bin/bun apps/mcp/stealth-request/server.ts
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5

# Environment
Environment="PUPPETEER_CACHE_DIR=/root/.cache/puppeteer"
Environment="NODE_ENV=production"

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stealth-request

[Install]
WantedBy=multi-user.target
