# apps/mcp-servers/browser-control/browser-control.service
[Unit]
Description=Browser Control Service for Alive Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/webalive/alive/apps/mcp-servers/browser-control
ExecStart=/usr/local/bin/bun src/server.ts
Restart=always
RestartSec=5
Environment="NODE_ENV=production"
EnvironmentFile=/root/webalive/alive/apps/mcp-servers/browser-control/.env
# Chromium + pool overhead
MemoryMax=1536M
CPUQuota=50%
# Playwright needs /dev/shm
PrivateDevices=no

[Install]
WantedBy=multi-user.target
