# apps/mcp-servers/stealth-request/stealth-request.service
[Unit]
Description=Stealth Request Web Scraping Service
After=network.target
Documentation=file:///root/webalive/alive/apps/mcp-servers/stealth-request/README.md

[Service]
Type=simple
User=root
WorkingDirectory=/root/webalive/alive
ExecStart=/root/.bun/bin/bun apps/mcp-servers/stealth-request/server.ts
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5

# Auto-install matching Chrome before start (prevents version mismatch after updates)
ExecStartPre=/usr/bin/npx puppeteer browsers install chrome

# Environment
Environment="PUPPETEER_CACHE_DIR=/root/.cache/puppeteer"
Environment="NODE_ENV=production"

# Resource limits
MemoryMax=1536M
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stealth-request

[Install]
WantedBy=multi-user.target
